# debian/postrm
#!/bin/sh

update-rc.d -f fake_ilo remove
